# datastore/supervisord.conf
[supervisord]
nodaemon=true

[program:datastore]
command=gcloud beta emulators datastore start --host-port gdatastore:8081